# board/link/userspace/etc/init/ec.conf
# Copyright (c) 2012 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description    "Overrides lightbar params to tweak appearance"
author         "chromium-os-dev@chromium.org"

start on startup

script
  PARAMS=/usr/share/ec/lightbar_params.txt
  [ -f "$PARAMS" ] && ectool lightbar params "$PARAMS" >/dev/null 2>&1
end script
